# debian/hd-idle.postinst
#!/bin/bash

if [ $(systemctl is-enabled hd-idle) != 'enabled' ];  then
  systemctl enable hd-idle
fi
systemctl daemon-reload
systemctl restart hd-idle
